# doraemon/系統操作/mac-address-changer.conf
[program:mac-address-changer.conf]
command = macchanger -m ac:1f:6b:62:8b:f0 eth0

autostart = true
autorestart = true
stopwaitsecs = 60
stderr_logfile = /var/log/supervisor/%(program_name)s.err.log
stdout_logfile = /var/log/supervisor/%(program_name)s.out.log
stderr_logfile_maxbytes = 1MB
stdout_logfile_maxbytes = 1MB
stderr_logfile_backups = 1
stdout_logfile_backups = 1

